# example-configs/supervisord.conf
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=./tmp/supervisor.sock

[supervisorctl]
serverurl=unix://./tmp/supervisor.sock

[supervisord]
logfile=./logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=./tmp/pids/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[program:sms_telnet_transport]
stdout_logfile = ./logs/%(program_name)s_%(process_num)s.log
stderr_logfile = ./logs/%(program_name)s_%(process_num)s.log
command = twistd -n --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
          start_worker
          --worker-class=vumi.transports.telnet.AddressedTelnetServerTransport
          --config=example-configs/sms_telnet_transport.yaml
          --vhost=/develop

[program:ussd_telnet_transport]
stdout_logfile = ./logs/%(program_name)s_%(process_num)s.log
stderr_logfile = ./logs/%(program_name)s_%(process_num)s.log
command = twistd -n --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
          start_worker
          --worker-class=vumi.transports.telnet.AddressedTelnetServerTransport
          --config=example-configs/ussd_telnet_transport.yaml
          --vhost=/develop

[program:wikipedia_dispatcher]
stdout_logfile = ./logs/%(program_name)s_%(process_num)s.log
stderr_logfile = ./logs/%(program_name)s_%(process_num)s.log
command = twistd -n --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
          start_worker
          --worker-class=vumi.dispatchers.endpoint_dispatchers.RoutingTableDispatcher
          --config=example-configs/wikipedia_dispatcher.yaml
          --vhost=/develop

[program:wikipedia_app]
stdout_logfile = ./logs/%(program_name)s_%(process_num)s.log
stderr_logfile = ./logs/%(program_name)s_%(process_num)s.log
command = twistd -n --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
          start_worker
          --worker-class=vumi_wikipedia.wikipedia.WikipediaWorker
          --config=example-configs/wikipedia_app.yaml
          --vhost=/develop
